# config/neo-node.service
[Unit]
Description=Neo N3 Rust node (neo-node)
Documentation=https://github.com/r3e-network/neo-rs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust User/Group and paths to your deployment. The binary is built with
#   cargo build --release -p neo-node --features wip
# and installed to /usr/local/bin/neo-node.
User=neo
Group=neo
ExecStart=/usr/local/bin/neo-node --config /etc/neo/neo_mainnet_node.toml
Restart=on-failure
RestartSec=5
# RocksDB data directory (also set [storage].data_dir / .path in the TOML).
StateDirectory=neo
WorkingDirectory=/var/lib/neo
Environment=RUST_LOG=info
# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
